# setup/webui_service.service
[Unit]
Description = Start ROS 2 Web UI
After = network.target

[Service]
Type=simple
User=%i
ExecStart=%i


[Install]
WantedBy = multi-user.target